# systemd-services/dnsupdate.service
# /etc/systemd/system/dnsupdate.service
[Unit]
Description=Update DNS records

[Service]
Type=oneshot
ExecStart=/home/pol/dnsupdate.sh
